# systemd/boot-sync.service
[Unit]
Description=Boot-time ROM sync (optional)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/kiosk-retropie/config.env
ExecStart=/usr/local/lib/kiosk-retropie/sync-roms.sh

[Install]
WantedBy=multi-user.target
