# systemd/emergency-retro-launch.service
[Unit]
Description=Emergency RetroPie launch listener (TTY)
After=multi-user.target

[Service]
Type=simple
EnvironmentFile=-/etc/kiosk-retropie/config.env
ExecStart=/usr/local/lib/kiosk-retropie/controller-listener-tty.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
