# systemd/ha-kiosk.service
[Unit]
Description=Home Assistant kiosk mode
After=network-online.target
Conflicts=retro-mode.service
Wants=network-online.target

# Fail open to RetroPie if HA kiosk repeatedly fails.
OnFailure=retro-ha-failover.service
OnFailureJobMode=replace
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
EnvironmentFile=-/etc/retro-ha/config.env
User=retropi
Group=retropi

# Create a logind session so rootless Xorg can acquire the seat.
PAMName=login
StandardInput=tty
TTYPath=/dev/tty7
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

ExecStart=/usr/local/lib/retro-ha/ha-kiosk.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
