# systemd/healthcheck.service
[Unit]
Description=Appliance healthcheck

[Service]
Type=oneshot
EnvironmentFile=-/etc/kiosk-retropie/config.env
ExecStart=/usr/local/lib/kiosk-retropie/healthcheck.sh
