# systemd/kiosk-retropie-install.service
[Unit]
Description=kiosk-retropie one-time installer
Wants=network-online.target
After=network-online.target

# Never re-run once successful.
ConditionPathExists=!/var/lib/kiosk-retropie/installed

# Avoid systemd start-rate limiting while retrying.
StartLimitIntervalSec=0

[Service]
Type=oneshot
EnvironmentFile=-/etc/kiosk-retropie/config.env

# Bootstrap is responsible for acquiring the installer and executing it.
ExecStart=/usr/local/lib/kiosk-retropie/bootstrap.sh

# Retry until network is available and installation succeeds.
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
