# systemd/kiosk-retropie-led-mqtt.service
[Unit]
Description=kiosk-retropie LED control via MQTT
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/kiosk-retropie/config.env
ExecStart=/usr/local/lib/kiosk-retropie/kiosk-retropie-led-mqtt.sh
Restart=on-failure
RestartSec=2

# Keep logs in journald (stdout/stderr).
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
