# systemd/kiosk.service
[Unit]
Description=Kiosk mode
After=network-online.target
Conflicts=retro-mode.service
Wants=network-online.target

# Fail open to RetroPie if kiosk repeatedly fails.
OnFailure=kiosk-retropie-failover.service
OnFailureJobMode=replace
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
EnvironmentFile=-/etc/kiosk-retropie/config.env
User=retropi
Group=retropi

# Create a logind session so rootless Xorg can acquire the seat.
PAMName=login
StandardInput=tty
TTYPath=/dev/tty7
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

ExecStart=/usr/local/lib/kiosk-retropie/kiosk.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
