# systemd/retro-ha-install.service
[Unit]
Description=retro-ha one-time installer
Wants=network-online.target
After=network-online.target

# Never re-run once successful.
ConditionPathExists=!/var/lib/retro-ha/installed

[Service]
Type=oneshot
EnvironmentFile=-/etc/retro-ha/config.env

# Bootstrap is responsible for acquiring the installer and executing it.
ExecStart=/usr/local/lib/retro-ha/bootstrap.sh

# Retry until network is available and installation succeeds.
Restart=on-failure
RestartSec=30

# Avoid systemd start-rate limiting while retrying.
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
