# systemd/save-backup.service
[Unit]
Description=Backup RetroPie saves/states to NFS (optional)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/kiosk-retropie/config.env
ExecStart=/usr/local/lib/kiosk-retropie/save-backup.sh
